# project/node.conf
# node

start on startup

script
	cd www
	node app.js
end script
